# debian/postinst
#!/bin/bash
set -e

WORKDIR="/opt/vyos-webui"
BACKEND_DIR="$WORKDIR/backend"
VENV_DIR="$BACKEND_DIR/venv"

fix_venv() {
    echo "Fixing venv for VyOS Python 3.11..."

    cd "$BACKEND_DIR"

    if [ -d "venv" ]; then
        # Update pyvenv.cfg (like deploy.sh)
        cat > venv/pyvenv.cfg << 'CFGEOF'
home = /usr/bin
include-system-site-packages = false
version = 3.11.2
CFGEOF
        # Fix symlinks
        cd venv/bin
        rm -f python python3 python3.12
        ln -sf /usr/bin/python3 python
        ln -sf /usr/bin/python3 python3
        if [ -x /usr/bin/python3.11 ]; then
            ln -sf /usr/bin/python3.11 python3.11
        fi
        # Fix shebangs in all scripts
        for f in *; do
            if [ -f "$f" ] && [ -x "$f" ]; then
                head -n 1 "$f" | grep -q "^#!" && sed -i '1s|^#!.*python.*|#!/usr/bin/python3|' "$f" 2>/dev/null || true
            fi
        done
        cd ../..
    fi

    # Create .env file
    if [ ! -f "$BACKEND_DIR/.env" ]; then
        cat > "$BACKEND_DIR/.env" << 'EOF'
VYOS_HOST=127.0.0.1
VYOS_PORT=22
VYOS_USERNAME=vyos
VYOS_PASSWORD=vyos
VYOS_TIMEOUT=30
EOF
    fi

    echo "venv fix complete."
}

case "$1" in
    configure)
        # Fix the venv (like deploy.sh)
        fix_venv

        # Set permissions - vyos user, users group
        chown -R vyos:users /opt/vyos-webui 2>/dev/null || true
        chmod +x /opt/vyos-webui/vyos-webui-wrapper
        chmod +x /opt/vyos-webui/start.sh
        chmod +x /opt/vyos-webui/stop.sh
        chmod +x /opt/vyos-webui/status.sh
        chmod 600 /etc/vyos-webui/config.env 2>/dev/null || true

        # Reload systemd
        systemctl daemon-reload || true

        echo ""
        echo "========================================="
        echo "  VyOS Web UI installed successfully!"
        echo "========================================="
        echo ""
        echo "方式一: 使用 systemd（推荐）"
        echo "  sudo systemctl start vyos-webui"
        echo "  sudo systemctl enable vyos-webui"
        echo ""
        echo "方式二: 使用脚本（与 deploy.sh 一致）"
        echo "  cd /opt/vyos-webui && ./start.sh"
        echo "  cd /opt/vyos-webui && ./stop.sh"
        echo "  cd /opt/vyos-webui && ./status.sh"
        echo ""
        echo "访问地址: http://<your-vyos-ip>:8000"
        echo "默认登录: vyos / vyos"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
esac

exit 0
