# debian/vyos-webui.service
[Unit]
Description=VyOS Web UI
After=network.target

[Service]
Type=simple
User=vyos
Group=vyos
WorkingDirectory=/opt/vyos-webui/backend
EnvironmentFile=/etc/vyos-webui/config.env
ExecStart=/opt/vyos-webui/vyos-webui-wrapper
Restart=always
RestartSec=10

# Security settings - relaxed to allow venv setup
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true
ReadWritePaths=/opt/vyos-webui /tmp

[Install]
WantedBy=multi-user.target
